# systemd/claude-telegram.service
[Unit]
Description=Claude Code Telegram bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/tdporter/Documents/Programming/claude-telegram
ExecStart=/usr/bin/npx tsx src/index.ts
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
